# battery-monitor.service
[Unit]
Description=Battery Monitor

[Service]
ExecStart=/path/to/your/batter-monitor
Restart=always

[Install]
WantedBy=multi-user.target
